# samples/amon.supervisor.conf
[program:odinsys]
command=/opt/amon/amon
directory=/opt/amon
user=vassilux
autostart=true
autorestart=true
